# deploy/kpm33b-config-sender.service
[Unit]
Description=KPM33B Config Sender (Configuration Flow)
After=network-online.target kpm33b-broker.service kpm33b-proxy.service
Wants=network-online.target
Requires=kpm33b-broker.service

[Service]
Type=simple
WorkingDirectory=/opt/kpm33b_proxy
ExecStart=/opt/kpm33b_proxy/venv/bin/python /opt/kpm33b_proxy/run_config_sender.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
