# deploy/kpm33b-broker.service
[Unit]
Description=KPM33B Internal MQTT Broker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kpm33b
Group=kpm33b
ExecStart=/usr/sbin/mosquitto -c /opt/kpm33b_proxy/scripts/kpm33b_broker.conf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
